# dev-vcs/gitlab-daemon/files/gitlab-daemon.service
[Unit]
Description=GitLab Daemon
After=network.target

[Service]
ExecStart=/usr/sbin/gitlab-daemon

[Install]
WantedBy=multi-user.target